# com.godoftrading.weekly-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.godoftrading.weekly-report</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ghp/Projects/GodOfTrading/tools/weekly_learning_report.py</string>
        <string>--telegram</string>
    </array>

    <!-- 매주 월요일 09:00 (로컬=KST) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/ghp/Projects/GodOfTrading</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/ghp</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/godoftrading_weekly_report.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/godoftrading_weekly_report_err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
