# com.cryptosignal.binance.fast.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cryptosignal.binance.fast</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ghp/Projects/GodOfTrading/main.py</string>
        <string>--fast-radar</string>
    </array>

    <key>StartInterval</key>
    <integer>180</integer>

    <key>WorkingDirectory</key>
    <string>/Users/ghp/Projects/GodOfTrading</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/ghp</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>AUTO_TRADE_EXCHANGE</key>
        <string>binance</string>
        <key>GOT_MARKET_DATA_EXCHANGE</key>
        <string>binance</string>
        <key>GOT_STATE_NAMESPACE</key>
        <string>binance</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/godoftrading_binance_fast.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/godoftrading_binance_fast_err.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
