# com.polymarket.paper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.polymarket.paper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ghp/Projects/GodOfTrading/polymarket_paper_bot.py</string>
    </array>

    <!-- 1분마다 BTC 5m Polymarket paper 후보 확인 + 정산. 실주문 없음. -->
    <key>StartInterval</key>
    <integer>60</integer>

    <key>WorkingDirectory</key>
    <string>/Users/ghp/Projects/GodOfTrading</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/ghp</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>POLYMARKET_PAPER_MILESTONE_REPORT_AT</key>
        <string>2026-07-10 06:00</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/polymarket_paper.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/polymarket_paper_err.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
