# tools/com.polymarket.wallet.radar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.polymarket.wallet.radar</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghp/Projects/GodOfTrading/.venv-poly/bin/python</string>
        <string>/Users/ghp/Projects/GodOfTrading/polymarket_wallet_pipeline.py</string>
        <string>--radar</string>
        <string>--json</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ghp/Projects/GodOfTrading</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/ghp</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>15</integer>
    <key>StandardOutPath</key>
    <string>/tmp/polymarket_wallet_radar.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/polymarket_wallet_radar_err.log</string>
</dict>
</plist>
